# deploy/roles/backend/templates/althea-backend.service.j2
[Unit]
Description=Althea Link Backend
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=althea-link-backend --https \
--cert-file=/etc/letsencrypt/live/{{ link_server_name }}/fullchain.pem \
--key-file=/etc/letsencrypt/live/{{ link_server_name }}/privkey.pem \
--pool-templates=36000 \
--pool-tokens="0x0412C7c846bb6b7DC462CF6B453f76D8440b2609,0x30dA8589BFa1E509A319489E014d384b87815D89,0x9676519d99E390A180Ab1445d5d857E3f6869065,0x7580bFE88Dd3d07947908FAE12d95872a260F2D8" \
--dex-contract="0x7580bFE88Dd3d07947908FAE12d95872a260F2D8" \
--multicall-contract="0x9726268F55d581d5F50c3853969010ACDCe7Cbff" \
--query-contract="0xB2F37Ba3CaDEc9aAf43BC716B1B86656be2d06Eb"

Restart=always
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
